# web/etc/supervisord.conf
[supervisord]
nodaemon=true
user=root
 
# [program:sshd]
# command=/usr/sbin/sshd -D
 
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=900
stdout_logfile= /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
username=www-data
autorestart=true

[program:php-fpm7.3]
command = /usr/local/sbin/php-fpm
autorestart=true
